# server/deb/debian/prerm
#!/bin/sh

systemctl disable gorram-server.service >/dev/null
systemctl stop gorram-server.service >/dev/null