# client/deb/debian/prerm
#!/bin/sh

set -e

case "$1" in 
    upgrade)
    ;;  
    remove)
        deb-systemd-helper disable gorram-client.service >/dev/null
        deb-systemd-invoke stop gorram-client.service >/dev/null
    ;;
    *)
    ;;
esac

exit 0