# client/deb/debian/prerm
#!/bin/sh

systemctl disable gorram-client.service >/dev/null
systemctl stop gorram-client.service >/dev/null